# com.storyflow.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.storyflow.server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/devserver/local/node-v22.14.0-darwin-arm64/bin/node</string>
		<string>/Users/devserver/Documents/GitHub/StoryFlow/server/index.js</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/devserver/Documents/GitHub/StoryFlow</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/devserver/Library/Logs/storyflow.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/devserver/Library/Logs/storyflow.log</string>
</dict>
</plist>
